# meterparser/rootfs/etc/services.d/meterparser/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs meterparser
# ==============================================================================

declare -a options

options+=(--token ${SUPERVISOR_TOKEN})
if bashio::debug; then
  options+=(--debug)
fi
bashio::log.info "Starting meterparser..."

exec python3 /src "${options[@]}"
